# root/etc/init.d/tvgate
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

BIN_PATH=/usr/bin/tvgate/TVGate
CONFIG_PATH=/etc/tvgate/config.yaml

start_service() {
	config_load tvgate
	local enabled proxy listen_port config_port
	
	config_get_bool enabled enabled 0
	[ "$enabled" -eq 0 ] && return 1
	
	config_get proxy proxy 'https://hk.gh-proxy.com/'
	config_get listen_port listen_port '8888'
	
	# Create config directory
	mkdir -p /etc/tvgate
	
	# Check if binary exists, if not try to download
	if [ ! -f "$BIN_PATH" ]; then
		echo "TVGate binary not found, trying to download..."
		/usr/bin/tvgate-download.sh
	fi
	
	# Start service if binary exists
	if [ -f "$BIN_PATH" ]; then
		procd_open_instance
		procd_set_param command "$BIN_PATH" "-config=$CONFIG_PATH"
		procd_set_param respawn
		procd_close_instance
	else
		echo "TVGate binary not found, service not started"
		return 1
	fi
}

stop_service() {
	echo "Stopping TVGate service"
}

reload_service() {
	echo "Reloading TVGate service"
	stop
	start
}

service_triggers() {
	procd_add_reload_trigger "tvgate"
}